# pomodoro-motd-build.service
[Unit]
Description=Pomodoro MOTD build (extract key passages from pending sources)

[Service]
Type=oneshot
ExecStart=%h/pomodoro/pomodoro motd build
# Needs OPENROUTER_API_KEY — set it here or drop it in an EnvironmentFile.
# Environment=OPENROUTER_API_KEY=...
