# pomodoro-motd-build.service
[Unit]
Description=Pomodoro MOTD build (extract key passages from pending sources)

[Service]
Type=oneshot
# Reads OPENROUTER_API_KEY from the repo .env (the CLI also loads it itself).
EnvironmentFile=-%h/pomodoro/.env
ExecStart=%h/pomodoro/pomodoro motd build
