# pomodoro-learning-nudge.service
[Unit]
Description=Pomodoro daily-learning nudge (notify if nothing logged today)

[Service]
Type=oneshot
ExecStart=%h/pomodoro/pomodoro learning nudge
