# pomodoro-motd-build.timer
[Unit]
Description=Weekly rebuild of the pomodoro message-of-the-day pool

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
